# example/supervisor/configs/test.ini
[program:test]
command=python /opt/services/test/main.py
directory=/opt/services/test
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/test.log
stderr_logfile=/var/log/supervisor/test.err.log
